# clima_sensor/cesga-sensor_opentsdb/publish_cesga_sensor.timer
[Unit]
Description=Run publish gpu information

[Timer]
OnBootSec=2
OnUnitInactiveSec=300
RandomizedDelaySec=10

[Install]
WantedBy=timers.target
